# btusb.service
[Unit]
Description=USB HID Proxy Service

[Service]
User=root
ExecStart=/home/pi/btproxy/btusb.sh

[Install]
WantedBy=multi-user.target
